# supervisor/supervisor-app.conf
# each program entry below is a separate terminal command.
# Each command MUST run in the foreground and stay running.
# If the command ever exits, the supervisor daemon will automatically run it again.
# Programs can be controlled like this: 'supervisorctl start fail2ban' 'supervisorctl stop fail2ban'
# supervisor writes program statuses in /var/log/supervisor

[supervisord]
nodaemon=true

[program:saslauthd_ldap]
startsecs=0
autostart=true
autorestart=true
command=/usr/sbin/saslauthd -d -a ldap -O /etc/saslauthd.conf
pidfile=/var/run/saslauthd/saslauthd.pid

[program:prosody]
startsecs=0
autostart=true
autorestart=true
user=prosody
command=prosody
pidfile=/var/run/prosody/prosody.pid
